# templates/rebuild-file-cache.timer.j2
[Unit]
Description=Hashicorp-proxy rebuild file cache timer

[Timer]
OnBootSec=5min
OnCalendar={{ hp_oncalendar }}

[Install]
WantedBy=timers.target